# scripts/pipeline/systemd/flow-review.service
[Unit]
Description=Flow Review — autonomous Claude cross-review trigger
Documentation=file://%h/workspace/workspace/web-of-trust/scripts/pipeline/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run from the repo root so claude picks up the project context
# (settings, allowed tools, MCP servers configured for this repo).
WorkingDirectory=%h/workspace/workspace/web-of-trust

# Wrap the script in flock so the lock is held for the full lifetime of
# the run. If a previous tick is still going, the second invocation
# exits immediately without doing anything (-n = nonblock, exit on
# contention; -E 0 = report success when contended so systemd does not
# count it as a failure).
ExecStart=/usr/bin/flock -n -E 0 /run/user/%U/flow-review.lock %h/workspace/workspace/web-of-trust/scripts/pipeline/flow-review.sh

# Soft caps on resource use.
TimeoutStartSec=15min
MemoryMax=2G

# Failure handling — if the script itself errors, systemd marks the
# unit failed; the script's own self-failure path also opens an issue.
Environment="PATH=/usr/local/bin:/usr/bin:/bin:%h/.nvm/versions/node/current/bin:%h/.local/bin"

[Install]
WantedBy=default.target
